# server/bin/charted.service
[Unit]
Description=You know, for Helm Charts?
Documentation=https://charts.noelware.org/docs
After=network-online.target

[Service]
Type=notify
TimeoutStartSec=0
ExecStart=/etc/noelware/charted/server/bin/charted-server --config=/etc/noelware/charted/server/config/charted.yml
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
